# build/dev/services/vite.conf
[program:vite]
command=npm run serve
user=app
priority=400
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
